# etc/profile.d/go.sh
export GOROOT=/usr/lib/go  #设置为go安装的路径，有些安装包会自动设置默认的goroot
export GOPATH=$HOME/works/go   #默认安装包的路径
export PATH=$PATH:$GOROOT/bin:$GOPATH/bin
