# config/autostart/say_weather.desktop
#!/usr/bin/env xdg-open
 [Desktop Entry]
Name=weather
Exec=sleep 30;sh /home/lunatic/bin/weather.sh
Type=Application
Terminal=false
